# infra/zeroclaw/systemd/zeroclaw-gateway.service
[Unit]
Description=ZeroClaw Gateway (ZHC-Nova)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/zhc-nova
EnvironmentFile=/opt/zhc-nova/.env
ExecStart=/usr/bin/env bash -lc 'set -euo pipefail; if [[ -z "${ZHC_ZEROCLAW_GATEWAY_CMD:-}" ]]; then echo "ERROR: set ZHC_ZEROCLAW_GATEWAY_CMD in .env" >&2; exit 1; fi; exec ${ZHC_ZEROCLAW_GATEWAY_CMD}'
Restart=always
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=20
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/zhc-nova/storage

[Install]
WantedBy=multi-user.target
